# axi-bot.service
[Unit]
Description=Axi Personal Assistant Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pride/coding-projects/personal-assistant
ExecStart=/usr/bin/uv run python /home/pride/coding-projects/personal-assistant/supervisor.py
Restart=on-failure
RestartSec=10
Environment="PATH=/usr/bin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=default.target
