# axi-bot.service
[Unit]
Description=Axi Personal Assistant Discord Bot
After=axi-procmux.service
Wants=axi-procmux.service
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/axi-assistant
EnvironmentFile=/home/ubuntu/axi-assistant/.env
ExecStart=/bin/bash -c 'exec "$AXI_RS_BINARY"'
Restart=always
RestartPreventExitStatus=0 SIGTERM
SuccessExitStatus=42
RestartSec=2
TimeoutStopSec=15
Environment="HOME=/home/ubuntu"
Environment="PATH=/home/ubuntu/.cargo/bin:/home/ubuntu/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=default.target
