# axi-procmux.service
[Unit]
Description=Axi Procmux (process multiplexer)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/axi-assistant
EnvironmentFile=%h/axi-assistant/.env
ExecStart=/bin/bash -c 'exec "$(dirname "$AXI_RS_BINARY")/procmux" %h/axi-assistant/.bridge.sock'
Restart=on-failure
RestartSec=1
TimeoutStopSec=10
Environment="HOME=%h"

[Install]
WantedBy=default.target
