# axi-bot.service
[Unit]
Description=Axi Personal Assistant Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/axi-assistant
ExecStart=/home/ubuntu/.local/bin/uv run python -m axi.supervisor
ExecReload=kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStopSec=15
Environment="PATH=/home/ubuntu/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=default.target
