# axi-rs/systemd/axi-test-procmux@.service
[Unit]
Description=Axi Test Procmux (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/axi-tests/%i
EnvironmentFile=/home/ubuntu/axi-tests/%i/.env
ExecStart=/bin/bash -c 'exec "$(dirname "$AXI_RS_BINARY")/procmux" /home/ubuntu/axi-tests/%i/.bridge.sock'
Restart=on-failure
RestartSec=1
TimeoutStopSec=10
Environment="HOME=/home/ubuntu"
